# supervisor/supervisord.d/nginx.ini
[program:nginx]
startsecs = 0
autorestart = false
command = nginx -c /usr/local/nginx/conf/nginx.conf
