# supervisor/supervisord.d/php.ini
[program:php]
startsecs = 0
autorestart = false
command = php-fpm start
